# src/crontab
0 * * * * /usr/bin/python /opt/datagrid/bin/slv_extractor.py >> /opt/datagrid/log/cron.log 2>&1
